# src/etc/hotplug.d/iface/95-aliddns
#!/bin/sh

[ -x /usr/bin/aliddns.sh ] || exit 0

case "${ACTION:-}" in
	ifup|ifupdate)
		/usr/bin/aliddns.sh --trigger "${ACTION}:${INTERFACE:-unknown}" "${INTERFACE:-}" &
		;;
esac